# ansible/roles/app/templates/puma.service.j2
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/{{ ansible_user }}/db_config
User={{ ansible_user }}
WorkingDirectory=/home/{{ ansible_user }}/reddit
ExecStart=/bin/bash -lc 'puma'
Restart=always

[Install]
WantedBy=multi-user.target
